# dockerfiles/supervisor-app.conf
[supervisord]
[program:workers]
directory = /app/BEacon
command = /usr/bin/uwsgi --ini /etc/uwsgi/apps-enabled/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx